# deploy/systemd/atta.service
[Unit]
Description=Automated Ticket Attribution (one run)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/atta/current
EnvironmentFile=/opt/atta/.env
ExecStart=/opt/atta/venv/bin/python -m app.cmd.main